# scripts/wenxian-autopush.service
[Unit]
Description=Wenxian auto-push watcher (skills/wenxian -> GitHub)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/chenyechao/.openclaw/workspace
ExecStart=/home/chenyechao/.openclaw/workspace/scripts/wenxian_autopush_watch.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
